# etc/profile.d/nvidia-discrete-only.sh
export GBM_BACKEND=nvidia-drm
export WLR_NO_HARDWARE_CURSORS=1
export __GLX_VENDOR_LIBRARY_NAME=nvidia
export VK_ICD_FILENAMES=/usr/share/vulkan/icd.d/nvidia_icd.json
export VK_LOADER_DRIVERS_SELECT=*nvidia*
export LIBVA_DRIVER_NAME=nvidia
export NVD_BACKEND=direct

#uncomment if you use VDPAU
#export VDPAU_DRIVER=nvidia

#UNCOMMENT BELOW if you use WINE with DXVK-NVAPI
#export DXVK_ENABLE_NVAPI=1
