# etc/profile.d/wayland.sh
export ELECTRON_OZONE_PLATFORM_HINT=auto
export EGL_PLATFORM=wayland
export SDL_VIDEODRIVER=wayland
export CLUTTER_BACKEND=wayland
#export PROTON_ENABLE_WAYLAND=1
